# snippets/systemd/simple.service
[Unit]
Description=<desc>

[Service]
Type=simple
ExecStart=<executable>

[Install]
WantedBy=multi-user.target
